# bin/rcup.cron
#!/usr/bin/env zsh

set -e

USER=adam

source ~/.zshenv

dsa
rcup
zrec
emk
