# bin/rcup.cron
#!/usr/bin/env zsh

set -e

USER=adam

source ~/.zshenv

dsa
echo "\n-----------------\n"
rcup
echo "\n-----------------\n"
zrec
echo "\n-----------------\n"
emk
